# systemd/lay-daemon.service
[Unit]
Description=lay — Caramba/Punto-style layout switcher daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/lay-daemon
# Для диагностики: добавь --debug-log или -v. Вывод попадёт в journal и может содержать набранный текст.
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
